# deploy/linux/marketsharp_queue_worker_event.service
[Unit]
Description=MarketSharp Queue Worker (Event-Driven)
After=network.target

[Service]
Type=simple
Restart=always
RestartSec=5
ExecStart=/usr/bin/python3 /path/to/spicer/queue_ui_poster.py --event-driven

[Install]
WantedBy=multi-user.target

