# deploy/linux/spicer-webhook-sync.service
[Unit]
Description=Sync CompanyCam webhook URL to current Cloudflared tunnel
After=network-online.target spicer-cloudflared.service
Wants=network-online.target

[Service]
Type=oneshot
User=rellis
WorkingDirectory=/home/rellis/spicer
EnvironmentFile=/home/rellis/spicer/.env
ExecStart=/home/rellis/spicer/deploy/linux/sync_companycam_webhook.sh

[Install]
WantedBy=multi-user.target